# auf_dem_pi_einzurichten/trigger-mon.path
[Unit]
Description="Monitor the /home/pi/dvc/trigger file for changes"

[Path]
PathModified=/home/pi/dvc/trigger
Unit=trigger-mon.service

[Install]
WantedBy=multi-user.target
